# br2rauc/board/raspberrypi/rootfs-overlay/usr/lib/systemd/system/usb-gadget.service
[Unit]
Description=Set up USB gadget

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/usb-gadget.sh
